# charts/docker-novnc-prusaslicer/docker/app/conf.d/prusaslicer.conf
[program:prusaslicer]
command=launch-prusaslicer
autorestart=true